# config/supervisord.d/celery.ini
[program:celery_flower]
command=/usr/bin/python manage.py celery flower
directory=/saas/projects/AnsibleWeb/
startretries=3

[program:celery_worker]
command=/usr/bin/python manage.py celery worker -l info
directory=/saas/projects/AnsibleWeb/
stdout_logfile=/saas/projects/AnsibleWeb/logs/celery/celery_output.log
stdout_logfile_maxbytes=200MB
stderr_logfile=/saas/projects/AnsibleWeb/logs/celery/celery_error.log
stderr_logfile_maxbytes=200MB
startretries=3
stopasgroup=true
